# sv/packetbeat/run
#!/bin/sh

#sv start elasticsearch || exit 1

exec 2>&1
exec packetbeat -e -c /etc/packetbeat/packetbeat.conf 
