# sv/packetbeat/run
#!/bin/bash

#sv start elasticsearch || exit 1

source /etc/envvars
env

cd /packetbeat

exec 2>&1
exec ./packetbeat -c /packetbeat/packetbeat.yml
